# deploy/systemd/spark-agent.user.service
[Unit]
Description=Spark Signals host telemetry agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/projects/spark-signals/deploy/runtime/agent.env
ExecStart=%h/projects/spark-signals/target/release/spark-agent --config %h/projects/spark-signals/deploy/example-config/agent.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
LockPersonality=yes
MemoryDenyWriteExecute=yes
ReadOnlyPaths=/proc
ReadOnlyPaths=/sys

[Install]
WantedBy=default.target
