# deploy/systemd/spark-agent.service
[Unit]
Description=Spark Signals host telemetry agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spark-signals-agent
Group=spark-signals-agent
EnvironmentFile=/etc/spark-signals/agent.env
ExecStart=/usr/local/bin/spark-agent --config /etc/spark-signals/agent.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
UMask=0027
StateDirectory=spark-signals
StateDirectoryMode=0750

NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
ReadOnlyPaths=/proc
ReadOnlyPaths=/sys

[Install]
WantedBy=multi-user.target
